# ops/launchd/com.pantheon.agy-content-publisher.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pantheon.agy-content-publisher</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>scripts.agy_content_publisher</string>
    <string>--repo-root</string>
    <string>__REPO_ROOT__</string>
    <string>--queue-root</string>
    <string>__QUEUE_ROOT__</string>
    <string>--state-root</string>
    <string>__REPO_ROOT__/.work/content-publisher</string>
    <string>--max-runs</string>
    <string>__MAX_RUNS__</string>
    <string>--include-rewrites</string>
    <string>--push</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>

  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/agy-content-publisher.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/agy-content-publisher.stderr.log</string>
</dict>
</plist>
